# crates/la-daemon/templates/dev.lazyagents.lad.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd user agent for the LazyAgents daemon (M4.1 / A3 / S1).

  Generated by `lad install --service launchd`. Hard rules from M4
  brief v1.1 addendum (A3):
    * `ProgramArguments[0]` is an absolute path resolved at install time
      (`which lad`, falling back to `env::current_exe()`).
    * `ProgramArguments[1..]` is `["--config", "<absolute config path>"]`.
    * Placeholders substituted at install time:
        {{ exec_path }}    absolute path to the `lad` binary
        {{ config_path }}  absolute path to config.toml
        {{ home_path }}    resolved $HOME
        {{ user }}         resolved $USER
        {{ stdout_path }}  $HOME/Library/Logs/lazyagents/lad.out.log
        {{ stderr_path }}  $HOME/Library/Logs/lazyagents/lad.err.log
    * NO tilde, $HOME, $XDG_*, or other shell variables are allowed
      anywhere in this file after substitution. launchd does not
      expand them.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.lazyagents.lad</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ exec_path }}</string>
        <string>start</string>
        <string>--config</string>
        <string>{{ config_path }}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LAZYAGENTS_MANAGED_BY</key>
        <string>launchd</string>
        <key>LAZYAGENTS_CONFIG</key>
        <string>{{ config_path }}</string>
        <key>HOME</key>
        <string>{{ home_path }}</string>
        <key>USER</key>
        <string>{{ user }}</string>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>{{ stdout_path }}</string>

    <key>StandardErrorPath</key>
    <string>{{ stderr_path }}</string>
</dict>
</plist>
