# crates/la-daemon/templates/lad-task.xml
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Windows Scheduled Task definition for the LazyAgents daemon
  (M4.1 / A8 / S1).

  Generated by `lad install --service windows-task`. Hard rules from
  M4 brief v1.1 addendum (A8) + R3 decision:
    * Triggers: LogonTrigger only — runs at this user's interactive
      logon. No SCM service path.
    * Principal: the installing user, NOT SYSTEM and NOT an admin
      account. RunLevel = LeastPrivilege.
    * Settings: AllowHardTerminate=false + RestartOnFailure at least
      one retry; AllowDemandStart=true so `schtasks /Run` works.
    * `Enabled=false` on first install. `lad install --service ... --enable`
      flips it on via `schtasks /Change /TN ... /ENABLE`.
    * Placeholders substituted at install time:
        {{ user_id }}     domain\username of the installing user
        {{ exec_path }}   absolute path to lad.exe
        {{ config_path }} absolute path to config.toml
        {{ working_dir }} %APPDATA%\lazyagents (resolved)
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <URI>\LazyAgents\lad</URI>
    <Description>LazyAgents daemon (lad) — runs at user logon. Installed by `lad install --service windows-task`. Only runs while the user is logged in.</Description>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <UserId>{{ user_id }}</UserId>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>{{ user_id }}</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>false</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <!-- A8 hard rule: Enabled=false on install; --enable flips it on. -->
    <Enabled>false</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>{{ exec_path }}</Command>
      <Arguments>start --config "{{ config_path }}"</Arguments>
      <WorkingDirectory>{{ working_dir }}</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
