# crates/la-daemon/templates/lad.service
# systemd user unit for the LazyAgents daemon (M4.1 / A2 / S1).
#
# Generated by `lad install --service systemd`. Placeholders substituted
# at install time:
#   {{ exec_path }}    absolute path to the `lad` binary
#   {{ config_path }}  absolute path to config.toml (resolve_config_path())
#
# Do NOT hand-edit a copy that was already written to disk — rerun
# `lad install --service systemd` so the next upgrade keeps the file in
# sync. The service stays disabled by default; pass `--enable` to mark
# it enabled and `--start` to start it.

[Unit]
Description=LazyAgents daemon (lad)
Documentation=https://github.com/TinyWhite1997/LazyAgents
After=default.target

[Service]
Type=simple
ExecStart={{ exec_path }} start --config {{ config_path }}
Restart=on-failure
RestartSec=5
# S1: tell the la bootstrap that the daemon is managed by a service
# supervisor so it skips the auto-daemonize path.
Environment=LAZYAGENTS_MANAGED_BY=systemd
Environment=LAZYAGENTS_CONFIG={{ config_path }}

[Install]
WantedBy=default.target
